# systemd/discord-timeblock.service.in
[Unit]
Description=Enforce local Discord time block

[Service]
Type=oneshot
ExecStart=@HOME@/.local/bin/discord-guard enforce
